# hortiradar/website/anomaly.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

HORTI=/home/rahiel/hortiradar

# m h dom mon dow user  command
15 * * * * rahiel cd $HORTI/hortiradar/website && chronic $HORTI/venv/bin/python ./prediction.py 2>&1 | telegram-send -g --stdin --pre
