# database/streamer-supervisor.conf
[program:hortiradar-streamer]
command=/home/rahiel/hortiradar/venv/bin/python ./streamer.py
directory=/home/rahiel/hortiradar/database
autostart=yes
user=rahiel
