# website/caching.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

HORTI=/home/rahiel/hortiradar

# m h dom mon dow user  command
0 * * * * rahiel cd $HORTI/website && $HORTI/venv/bin/python ./refresh_cache.py
