# hortiradar/database/clean.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

HORTI=/home/rahiel/hortiradar

# m h dom mon dow user  command
30 2 * * * rahiel cd $HORTI/hortiradar/database && $HORTI/venv/bin/python ./clean.py
