# hortiradar/database/restart-workers.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

HORTI=/home/rahiel/hortiradar

# m h dom mon dow user  command
00 4 * * * root cd $HORTI/hortiradar/database && $HORTI/venv/bin/python ./restart_workers.py
